# deploy/cord.service
[Unit]
Description=cord message relay
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/cord
ExecStart=/usr/bin/node src/index.js
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
